# config/com.wooki.oss-radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wooki.oss-radar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wooki/project/toy/oss-radar/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/wooki/project/toy/oss-radar</string>

    <!-- Every Monday at 09:00 KST -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/wooki/project/toy/oss-radar/logs/cron.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wooki/project/toy/oss-radar/logs/cron.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/wooki</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/wooki/.local/bin</string>
    </dict>
</dict>
</plist>
